# wakeup.service
[Unit]
Description=Configure RTC wakeup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wakeup.sh
